# .config/systemd/user/color-scheme-watcher.service
[Unit]
Description=Monitor color-scheme changes

[Service]
ExecStart=%h/.local/lib/color-scheme-watcher
Restart=always
