# .config/systemd/user/hyprlock.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecStart=sh -c 'niri msg action do-screen-transition && hyprlock'
KillSignal=USR1
RestartKillSignal=TERM
Slice=session-graphical.slice
