# .config/systemd/user/playerctld.service
[Unit]
After=graphical-session.target

[Service]
Type=exec
ExecStart=playerctld
Slice=app-graphical.slice
