# .config/systemd/user/swaybg-restart.path
[Path]
PathModified=%E/background

[Install]
WantedBy=graphical-session.target
