# .config/systemd/user/xwayland-satellite.service.d/override.conf
[Service]
Slice=background-graphical.slice
ExecStart=
ExecStart=/usr/bin/xwayland-satellite :0
